# apps/api/supervisord.conf
; supervisord.conf
; API(uvicorn) + Celery Worker + Celery Beat を1コンテナで管理
; Render Disk を共有することでファイル共有が可能になる

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ─────────────────────────────────────────
; 1. FastAPI (Uvicorn)
; ─────────────────────────────────────────
[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_PORT)s --workers 1 --log-level %(ENV_LOG_LEVEL)s
directory=/app
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

; ─────────────────────────────────────────
; 2. Celery Worker（音声・動画・アップロード処理）
; ─────────────────────────────────────────
[program:celery_worker]
command=celery -A app.jobs.celery_app worker --loglevel=info --concurrency=2 -Q youtube,ai,video,upload,default
directory=/app
autostart=true
autorestart=true
priority=20
; Uvicorn より後に起動（DB接続確立を待つ）
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

; ─────────────────────────────────────────
; 3. Celery Beat（週次スケジューラー）
; ─────────────────────────────────────────
[program:celery_beat]
command=celery -A app.jobs.celery_app beat --loglevel=info --schedule=/tmp/celerybeat-schedule
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"
